# freebsd/files/usr/local/etc/rc.d/camcontrol.local
#!/bin/sh

. /etc/rc.subr

sudo camcontrol apm /dev/ada0
sudo camcontrol aam /dev/ada0
sudo camcontrol idle /dev/ada0 -t 0
